# dot_config/autostart/bitwarden.desktop
[Desktop Entry]
Type=Application
Name=Bitwarden
Exec=bitwarden-desktop
StartupNotify=false
Terminal=false
