# dot_config/autostart/seafile.desktop
[Desktop Entry]
Type=Application
Name=Seafile
Exec=seafile-applet
StartupNotify=false
Terminal=false
